# systemd/llama-70b.service.example
[Unit]
Description=llama.cpp server for future 70B-class model
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=REPLACE_WITH_LOCAL_USER
WorkingDirectory=REPLACE_WITH_REPO_ROOT
Environment=CUDA_VISIBLE_DEVICES=0,1
Environment=RUN_FOREGROUND=1
ExecStart=REPLACE_WITH_REPO_ROOT/scripts/start-70b.sh
ExecStop=REPLACE_WITH_REPO_ROOT/scripts/stop-70b.sh
Restart=on-failure
RestartSec=5
TimeoutStartSec=180
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
