# systemd/llama-gemma4-31b.service
[Unit]
Description=llama.cpp server for Gemma 4 31B IT
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=REPLACE_WITH_LOCAL_USER
WorkingDirectory=REPLACE_WITH_REPO_ROOT
Environment=CUDA_VISIBLE_DEVICES=1
Environment=RUN_FOREGROUND=1
ExecStart=REPLACE_WITH_REPO_ROOT/scripts/start-gemma4-31b.sh
ExecStop=REPLACE_WITH_REPO_ROOT/scripts/stop-gemma4-31b.sh
Restart=on-failure
RestartSec=5
TimeoutStartSec=120
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
