# systemd/llama-qwen36.service
[Unit]
Description=llama.cpp server for Qwen3.6-35B-A3B
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kzkchd
WorkingDirectory=/home/kzkchd/llm
Environment=CUDA_VISIBLE_DEVICES=0
Environment=RUN_FOREGROUND=1
ExecStart=/home/kzkchd/llm/scripts/start-qwen36.sh
ExecStop=/home/kzkchd/llm/scripts/stop-qwen36.sh
Restart=on-failure
RestartSec=5
TimeoutStartSec=120
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
