# systemd/valkey-local.service
[Unit]
Description=Temporary local Valkey for local LLM orchestration on inference host
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=REPLACE_WITH_LOCAL_USER
WorkingDirectory=REPLACE_WITH_REPO_ROOT
ExecStart=REPLACE_WITH_REPO_ROOT/scripts/start-valkey.sh
ExecStop=REPLACE_WITH_REPO_ROOT/scripts/stop-valkey.sh
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
